# itf.service
[Unit]
Description=ITF
[Service]
ExecStart=/usr/sbin/itf -q
[Install]
WantedBy=multi-user.target